# pipeline.service
[Unit]
Description=Pipeline Queue Manager Daemon
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=pipeline
Group=pipeline
WorkingDirectory=/home/pipeline/pipeline
Environment="PATH=/home/pipeline/.conda/envs/pipeline/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/pipeline/pipeline/.env
ExecStart=/home/pipeline/.conda/envs/pipeline/bin/python3 /home/pipeline/pipeline/pipeline/bin/queue_daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pipeline-queue

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
